# docs/platform/neuralock/05-server-setup/assets/neuralock.service
[Unit]
Description=Neuralock Server - Distributed Key Management Node
Documentation=https://docs.neuralock.io
After=network.target redis.service
Wants=redis.service

[Service]
Type=notify
User=neuralock
Group=neuralock
WorkingDirectory=/opt/neuralock
Environment="PATH=/opt/neuralock/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="NEURALOCK_ENV=production"
Environment="PYTHONUNBUFFERED=1"

# Service startup
ExecStartPre=/opt/neuralock/venv/bin/python scripts/health_check.py --pre-start
ExecStart=/opt/neuralock/venv/bin/python -m uvicorn src.main:app \
    --host 0.0.0.0 \
    --port 8000 \
    --workers 1 \
    --log-config config/logging.yaml

# Graceful shutdown
ExecStop=/bin/kill -TERM $MAINPID
TimeoutStopSec=30

# Service restart policy
Restart=on-failure
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/neuralock/data /opt/neuralock/logs
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
MemoryDenyWriteExecute=true
LockPersonality=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=2G
CPUQuota=100%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=neuralock-server

[Install]
WantedBy=multi-user.target